# etc/crontabs/root
# This is UTC hour
0 20 * * * /root/TP4OpenW/sync_list.sh
0 21 * * * sleep 70 && touch /etc/banner && reboot
