# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

sleep 10
touch /tmp/system-startup

ipset create gfwlist hash:ip counters timeout 1200
/root/TP4OpenW/update_conf.sh /root/TP4OpenW/lists/list
/etc/init.d/dnsmasq restart

# !!! Modify: [60000][60002][61002][UPSTREAM_IP]
autossh -M 0 -f -o ServerAliveInterval=30 -o ServerAliveCountMax=2 -C -L 0.0.0.0:49080:localhost:49080 -N root@[UPSTREAM_IP]
autossh -M 0 -o ServerAliveInterval=30 -o ServerAliveCountMax=2 -f -C -R 0.0.0.0:61002:localhost:22 -N root@[UPSTREAM_IP]
autossh -M 0 -o ServerAliveInterval=30 -o ServerAliveCountMax=2 -f -C -L 5353:localhost:5353 -N root@[UPSTREAM_IP]
socat udp4-listen:55,reuseaddr,fork tcp:localhost:5353 &

iptables -t nat -N REDSOCKS
iptables -t nat -A REDSOCKS -d 0.0.0.0/8 -j RETURN
iptables -t nat -A REDSOCKS -d 10.0.0.0/8 -j RETURN
iptables -t nat -A REDSOCKS -d 127.0.0.0/8 -j RETURN
iptables -t nat -A REDSOCKS -d 169.254.0.0/16 -j RETURN
iptables -t nat -A REDSOCKS -d 172.16.0.0/12 -j RETURN
iptables -t nat -A REDSOCKS -d 192.168.0.0/16 -j RETURN
iptables -t nat -A REDSOCKS -d 224.0.0.0/4 -j RETURN
iptables -t nat -A REDSOCKS -d 240.0.0.0/4 -j RETURN
iptables -t nat -A REDSOCKS -p tcp -m set --match-set gfwlist dst -j REDIRECT --to-ports 12345
iptables -t nat -A PREROUTING -p tcp -j REDSOCKS

exit 0
